# docker/rootfs/etc/s6-overlay/s6-rc.d/init-010-update/up
/etc/s6-overlay/s6-rc.d/init-010-update/run